# Library/LaunchAgents/name.iandunn.update-checkouts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
	  <key>Label</key>
	  <string>com.alvin.crontabtest</string>

	  <key>ProgramArguments</key>
	  <array>
	    <string>/Users/iandunn/dotfiles/bin/update-checkouts.sh</string>
	  </array>

	  <key>Nice</key>
	  <integer>1</integer>

	  <key>StartInterval</key>
	  <integer>3600</integer>

	  <key>RunAtLoad</key>
	  <true />

	  <key>StandardErrorPath</key>
	  <string>/tmp/name.iandunn.update-checkouts-errors.log</string>

	  <key>StandardOutPath</key>
	  <string>/tmp/name.iandunn.update-checkouts-messages.log</string>
	</dict>
</plist>
